# server/friday-integration.service
[Unit]
Description=Friday Integration Server for OpenClaw
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /Users/shreyansh/Documents/ESP\ 32\ INFRA/server/friday_integration_server.py --port 41263
Restart=always
RestartSec=5
Environment="PYTHONUNBUFFERED=1"

# Run as user (not root)
User=shreyansh
Group=shreyansh

# Working directory
WorkingDirectory=/Users/shreyansh/Documents/ESP\ 32\ INFRA/server

# Create necessary directories
ExecStartPre=/bin/mkdir -p /Users/shreyansh/.friday/recordings /Users/shreyansh/.friday/logs

[Install]
WantedBy=default.target
